# agent.service
[Unit]
Description=Full Access AI Agent Server
After=network.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/ai-agent-server
ExecStart=/usr/bin/env python3 /home/user/ai-agent-server/agent_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
