# agent.service
[Unit]
Description=AI Agent Server (anonysec/ai-agent-server)
After=network.target
[Service]
Type=simple
User=$(whoami)
WorkingDirectory=$(pwd)
ExecStart=/usr/bin/env python3 $(pwd)/agent_server.py
Restart=always
RestartSec=5
[Install]
WantedBy=multi-user.target
